# buildroot/board/otter/ottercast-S3/skeleton/etc/systemd/system/pulseaudio.service.d/override.conf
[Service]
User=otter
Group=otter
ExecStart=
ExecStart=/usr/bin/pulseaudio --daemonize=no
