# buildroot/package/ifupdown_hotplug/src/80-ifupdown.rules
# Handle allow-hotplug interfaces
SUBSYSTEM=="net", ACTION=="add|remove", RUN+="ifupdown-hotplug"
